# root/etc/services.d/isc-dhcp-server/run
#!/usr/bin/with-contenv bash

# Drop privileges and run dhcpd
dhcpd -user abc -group abc -d -f -4 -cf /config/dhcpd.conf -lf /data/dhcpd.leases -pf /run/dhcpd.pid ${DHCP_IF}



