# root/etc/services.d/isc-dhcp-server/run
#!/usr/bin/with-contenv bash

# Drop privileges and run dhcpd
dhcpd -user abc -group abc -d -f -4 --no-pid -cf /config/dhcpd.conf -lf /config/dhcpd.leases ${DHCP_IF}



